# docker/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/2 * * * * . /etc/environment; cd /app && /usr/local/bin/python main.py >> /var/log/cron.log 2>&1
